# pr-review/systemd/pr-review.timer
[Unit]
Description=Run PR Review Bot at 10:00 and 15:00

[Timer]
OnCalendar=*-*-* 10:00:00
OnCalendar=*-*-* 15:00:00

[Install]
WantedBy=timers.target
